# supervisor_development.ini
[program:jamiecurle]
command=/Users/jcurle/.virtualenvs/jamiecurle/bin/python /Users/jcurle/sites/jamiecurle/manage.py run_gunicorn -b 127.0.0.1:9001 --setting=settings_staging
directory=/Users/jcurle/Sites/jamiecurle
user=jcurle
autostart=true
autorestart=true
redirect_stderr=true

